# data/rc.local
#!/bin/bash

sudo -i -u ballsbot autossh -fnNT ballsbot-gw
